# rootfs/etc/services.d/app/run
#!/bin/bash
# ==============================================================================
# Starts app
# ==============================================================================
python3 -m ovoenergy version
